# supervisord.conf
[supervisord]
nodaemon=true

[program:crawler]
command=/usr/share/nginx/html/parking_crawler/permanent_crawl.sh

[program:nginx]
command=nginx -g daemon off